# resources/upstart/dexter-celery.conf
description "dexter-celery"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid mma
setgid mma
chdir /home/mma/mma-dexter

exec /home/mma/mma-dexter/production-celery.sh
