# resources/upstart/dexter.conf
description "dexter"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid ubuntu
setgid ubuntu
chdir /home/ubuntu/mma-dexter

exec /home/ubuntu/mma-dexter/bin/run-dexter.sh
